# scripts/datavault-server.service
[Unit]
Description=datavault Server
After=network-online.target zfs.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/datavault-server --config /etc/datavault/server/config.yaml
Restart=always
RestartSec=10
PrivateTmp=yes
# Keep this unit compatible with systemd 219 (for example, CentOS 7).
# ProtectSystem=strict and ReadWritePaths were introduced later and are not
# parsed by that baseline.
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
